# packages/sx-sf-hdparm-sdb-idle/hdparm.service
[Unit]
Description=hdparm sleep sda

[Service]
Type=oneshot
ExecStart=/usr/bin/hdparm -q -S 2 -y /dev/sda

[Install]
WantedBy=multi-user.target

