# packages/sx-sf-hdparm-sdb-idle/sx-sf-hdparm-sdb-idle.install
post_install() {
  systemctl enable hdparm
  systemctl enable hdparm-apm
  systemctl start hdparm
}

post_remove() {
  systemctl disable hdparm
  systemctl disable hdparm-apm
}
